# service/memcached.d/run
#!/bin/sh
MEMCACHED=/usr/local/memcached/bin/memcached
$MEMCACHED\
 -u www\
 -m 64\
 -l 192.168.1.121\
 -p 11211\
 -I 4m\
 2>/dev/stdout
